# root/etc/cron.d/backup-config
# Execute configuration backup if enabled
10 2 * * *   root /sbin/e-smith/backup-config

