# root/etc/cron.d/backup-config
# Execute configuration backup if enabled
15 0 * * *   root /usr/libexec/nethserver/backup-config-cron

